# void-files/etc/cron.daily/xbps-cache-clean
#!/bin/bash

# Automatically prune obsolete package files from the cache to save space
xbps-remove -O
